# config/rtail/ouroboros-err.conf
[program:rtail-ouroboros-err]
command=/bin/bash -c "tail -F /var/log/ouroboros.err.log | /usr/local/bin/rtail --id ouroboros-error --mute"
user = ouroboros
autostart=true
autorestart=true
