# config/rtail/rtail-server.conf
[program:rtail-server]
command=/usr/local/bin/rtail-server
user = ouroboros
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true